# data/system/cron-principal-backup-2026-06-05-12-12/BridgeCollector2h.modified.xml
<?xml version="1.0" encoding="UTF-16"?><Task version="1.3" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task"><RegistrationInfo><Description>桥接采集 - 12:00-16:00 数据真空期填补 (每 2h 触发)</Description><URI>\BridgeCollector2h</URI></RegistrationInfo><Principals><Principal id="Author" LogonType="ServiceAccount" UserId="SYSTEM" RunLevel="Highest"><UserId>S-1-5-21-904784874-1342172122-3530307549-500</UserId><LogonType>InteractiveToken</LogonType></Principal></Principals><Settings><DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries><StopIfGoingOnBatteries>false</StopIfGoingOnBatteries><ExecutionTimeLimit>PT5M</ExecutionTimeLimit><MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy><StartWhenAvailable>true</StartWhenAvailable><IdleSettings><StopOnIdleEnd>true</StopOnIdleEnd><RestartOnIdle>false</RestartOnIdle></IdleSettings><UseUnifiedSchedulingEngine>true</UseUnifiedSchedulingEngine></Settings><Triggers><TimeTrigger><StartBoundary>2026-06-05T12:00:00+08:00</StartBoundary><Repetition><Interval>PT2H</Interval><Duration>P3650D</Duration><StopAtDurationEnd>true</StopAtDurationEnd></Repetition></TimeTrigger></Triggers><Actions Context="Author"><Exec><Command>powershell.exe</Command><Arguments>-ExecutionPolicy Bypass -File C:\Users\Administrator\clawd\agents\workspace-gid\scripts\bridge-2h-cron.ps1</Arguments><WorkingDirectory>C:\Users\Administrator\clawd\agents\workspace-gid</WorkingDirectory></Exec></Actions></Task>
